# files/ttl-override.service
[Unit]
Description=TTL Override
After=ql-netd.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/etc/initscripts/ttl-override start
User=root

[Install]
WantedBy=multi-user.target
